# docker-image/scripts/profile.d/docker.sh
function docker-select() {
    export DOCKER_HOST="tcp://${1}.node.in.dawn:2376"
}
